# myproxy/myproxy.cron.in
#!/bin/sh

# The MyProxy cron script deletes users proxies that have expired.

# Set GLOBUS_LOCATION as appropriate for your installation.
GLOBUS_LOCATION="/usr/local/globus"
export GLOBUS_LOCATION

. ${GLOBUS_LOCATION}/libexec/globus-script-initializer

exec myproxy-admin-query -e 0 -r
