# myproxy.cron
#!/bin/sh

# This cron script deletes invalid MyProxy credentials
# (expired, revoked, etc.).

# Set GLOBUS_LOCATION as appropriate for your installation.
GLOBUS_LOCATION="/usr/grid"
export GLOBUS_LOCATION

. ${GLOBUS_LOCATION}/libexec/globus-script-initializer

${GLOBUS_LOCATION}/sbin/myproxy-admin-query -i -r 2>&1 | logger -t myproxy.cron -p cron.info

exit 0
