# myproxy.cron
#!/bin/sh

# The MyProxy cron script deletes users proxies that have expired.

# Set GLOBUS_LOCATION as appropriate for your installation.
GLOBUS_LOCATION="/usr/grid"
export GLOBUS_LOCATION

. ${GLOBUS_LOCATION}/libexec/globus-script-initializer

${GLOBUS_LOCATION}/sbin/myproxy-admin-query -e 0 -r 2>&1 | logger -t myproxy.cron -p cron.info

exit 0
