# gram/jobmanager/auditing/source/globus-gram-audit.cron.in
#! /bin/sh

@GLOBUS_SCRIPT_INITIALIZER@

@sbindir@/globus-gram-audit --quiet

exit 0
