# gram/jobmanager/auditing/source/globus-gram-audit.cron.in
#! /bin/sh

@GLOBUS_SCRIPT_INITIALIZER@

if [ -x @sbindir@/globus-gram-audit ] ; then
@sbindir@/globus-gram-audit --quiet
fi

exit 0
