# packaging/debian/globus-gatekeeper/debian/globus-gatekeeper.prerm
#! /bin/sh

if which invoke-rc.d >/dev/null 2>&1; then
	invoke-rc.d globus-gatekeeper stop
else
	/etc/init.d/globus-gatekeeper stop
fi
