# packaging/debian/globus-gatekeeper/debian/globus-gatekeeper.postinst
#! /bin/sh

if which update-rc.d  > /dev/null 2>&1; then
	update-rc.d globus-gatekeeper defaults
fi

if which invoke-rc.d >/dev/null 2>&1; then
	invoke-rc.d globus-gatekeeper start
else
	/etc/init.d/globus-gatekeeper start
fi
