# packaging/debian/globus-scheduler-event-generator/debian/globus-scheduler-event-generator-progs.prerm
#! /bin/sh

if which invoke-rc.d >/dev/null 2>&1; then
	invoke-rc.d globus-scheduler-event-generator stop
else
	/etc/init.d/globus-scheduler-event-generator stop
fi
