# packaging/debian/globus-scheduler-event-generator/debian/globus-scheduler-event-generator-progs.postrm
#! /bin/sh

if [ "$1" = purge ]; then
    if which update-rc.d > /dev/null 2>&1; then
    	update-rc.d globus-scheduler-event-generator remove
    fi
fi
