# packaging/debian/globus-gridftp-server/debian/globus-gridftp-server-progs.postinst
#! /bin/sh

if which update-rc.d  > /dev/null 2>&1; then
	update-rc.d globus-gridftp-server defaults
	update-rc.d globus-gridftp-sshftp defaults
fi

if which invoke-rc.d >/dev/null 2>&1; then
	invoke-rc.d globus-gridftp-server start
	invoke-rc.d globus-gridftp-sshftp start
else
	/etc/init.d/globus-gridftp-server start
	/etc/init.d/globus-gridftp-sshftp start
fi
