# packaging/debian/globus-gridftp-server/debian/globus-gridftp-server-progs.postrm
#! /bin/sh

if [ "$1" = purge ]; then
    if which update-rc.d > /dev/null 2>&1; then
    	update-rc.d globus-gridftp-server remove
    	update-rc.d globus-gridftp-sshftp remove
    fi
fi
