# packaging/debian/globus-gram-job-manager-lsf/debian/globus-gram-job-manager-lsf-setup-seg.prerm
#! /bin/sh

case "$1" in
    remove)
        globus-gatekeeper-admin -d jobmanager-lsf-seg > /dev/null 2>&1 || :
        globus-scheduler-event-generator-admin -d lsf > /dev/null 2>&1 || :

        if which invoke-rc.d >/dev/null 2>&1; then
            invoke-rc.d globus-scheduler-event-generator stop lsf || :
        elif [ -x /etc/init.d/globus-scheduler-event-generator ]; then
            /etc/init.d/globus-scheduler-event-generator stop lsf || :
        fi
        ;;
    upgrade)
        if which invoke-rc.d >/dev/null 2>&1; then
            invoke-rc.d globus-scheduler-event-generator stop lsf || :
        elif [ -x /etc/init.d/globus-scheduler-event-generator ]; then
            /etc/init.d/globus-scheduler-event-generator stop lsf || :
        fi
        ;;
    deconfigure)
        globus-gatekeeper-admin -d jobmanager-lsf-seg > /dev/null 2>&1 || :
        ;;
    failed-upgrade)
        :
        ;;
esac
